# config/supervisord.conf
[supervisord] nodaemon=true

[program:cups]
command=/usr/sbin/cupsd -f
redirect_stderr=true
user=print

[program:gcp]
command=gcp-cups-connector --log-to-console
